# assets/noctalia-polkit.service.in
[Unit]
Description=Noctalia Polkit Authentication Agent
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
ExecStartPre=/bin/sh -c 'mkdir -p "$${XDG_DATA_HOME:-$$HOME/.local/share}/dbus-1/services" && cp -n @DATADIR@/noctalia-polkit/org.gnome.keyring.SystemPrompter.service "$${XDG_DATA_HOME:-$$HOME/.local/share}/dbus-1/services/"'
ExecStart=@LIBEXECDIR@/noctalia-polkit
Slice=session.slice
TimeoutStopSec=5sec
Restart=on-failure

[Install]
WantedBy=graphical-session.target
